# pi-kiosk/carbonpanel-kiosk.service
[Unit]
Description=CarbonPanel SPI display kiosk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__DIR__
ExecStart=__DIR__/.venv/bin/python __DIR__/kiosk.py __DIR__/config.json
Restart=always
RestartSec=3
User=__USER__

[Install]
WantedBy=multi-user.target
